# 2.0.1/vicidial-installer/systemd/services/vicidial-db-optimize.service
[Unit]
Description=VICIDIAL Database Optimization

[Service]
Type=oneshot
ExecStart=/usr/bin/perl /usr/share/astguiclient/AST_DB_optimize.pl
